# images/raspivirt-qemu+openwrt/setupfiles/update-motd-ip.path
[Unit]
Description=Monitor network changes to update MOTD
After=network.target

[Path]
# Watch for network state changes
PathModified=/run/systemd/netif/links/
PathModified=/run/systemd/netif/leases/

[Install]
WantedBy=multi-user.target
